# algo/ghost-vwap-algo.service
[Unit]
Description=Ghost Alpha VWAP Reclaim Algo Trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/mphinance/algo
EnvironmentFile=/home/mphinance/algo/secrets.env
ExecStart=/usr/bin/python3 -u /home/mphinance/algo/ghost_vwap_algo.py --live
Restart=on-failure
RestartSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ghost-vwap-algo

[Install]
WantedBy=multi-user.target
